# overlay/etc/init.d/alphasound-features
#!/sbin/openrc-run
# Reads ALPHASOUND_* feature flags from alphasound.txt and disables
# unwanted services BEFORE the default runlevel starts them. Net effect:
# turning Bluetooth off in the web UI saves the few seconds bluetoothd
# and bluealsad take to initialise — actual boot speedup, not just a
# config flag.
#
# Runlevel symlinks live in tmpfs (the apkovl is extracted fresh on
# every boot), so any deletion we do here is automatically undone on
# the next boot. Each boot starts from a clean baseline and re-applies
# whatever the current alphasound.txt says.

description="Apply feature toggles from alphasound.txt"

CONFIG=/media/mmcblk0p1/alphasound.txt

depend() {
    need localmount
}

start() {
    [ ! -f "$CONFIG" ] && return 0

    local ssh bt mode
    ssh=$(grep '^ALPHASOUND_SSHD=' "$CONFIG" 2>/dev/null | tail -1 | cut -d= -f2-)
    bt=$(grep '^ALPHASOUND_BLUETOOTH=' "$CONFIG" 2>/dev/null | tail -1 | cut -d= -f2-)
    mode=$(grep '^ALPHASOUND_MODE=' "$CONFIG" 2>/dev/null | tail -1 | cut -d= -f2-)
    [ -z "$ssh" ] && ssh=yes
    [ -z "$bt" ] && bt=yes
    [ -z "$mode" ] && mode=RUN

    if [ "$ssh" != "yes" ]; then
        ebegin "disabling sshd"
        rc-update del sshd default 2>/dev/null
        eend 0
    fi

    if [ "$bt" != "yes" ]; then
        ebegin "disabling Bluetooth"
        rc-update del bluetooth default 2>/dev/null
        rc-update del bluealsa default 2>/dev/null
        eend 0
    fi

    # mDNS is only useful in client mode (so you can ssh alphasound.local
    # or Roon can discover the device). When we're the AP in RUN mode
    # nobody's advertising anything interesting; skip avahi-daemon to
    # save ~1s of boot.
    if [ "$mode" = "RUN" ]; then
        ebegin "disabling avahi-daemon (AP mode)"
        rc-update del avahi-daemon default 2>/dev/null
        eend 0
    fi
}
