# overlay/etc/init.d/alphasound-clock
#!/sbin/openrc-run
# Seed the system clock from the apkovl's build timestamp.
#
# Pi Zero 2 W has no battery-backed RTC, so the kernel boots with the
# clock at epoch 1970. Every file we extracted from the apkovl then
# looks like it's from the far future, OpenRC logs "clock skew detected"
# and pessimistically regenerates its deptree, and every boot-runlevel
# service scribbles "adjusting mtime" lines to the console. Setting the
# clock to /etc/alphasound-version's mtime (stamped at apkovl build
# time) gives us a plausible-past time that's monotonically newer than
# anything on the tmpfs rootfs, which silences the warnings and saves
# a second or two of deptree churn. Real NTP corrects it properly once
# we're in client mode.

description="Seed system clock from apkovl build time"

depend() {
    before hwclock sysctl
}

start() {
    # /etc/alphasound-build-time is written at apkovl build time with a
    # plain ISO-format UTC timestamp ("YYYY-MM-DD HH:MM:SS"), which is
    # a form every busybox `date -s` accepts. We read the *contents*
    # rather than the file's mtime because mtimes across the rootfs
    # are intentionally backdated to epoch 0 at build time to silence
    # OpenRC's 'clock skew detected' warnings — so file mtimes can't
    # be trusted as a time source here.
    [ -f /etc/alphasound-build-time ] || return 0
    fmt=$(cat /etc/alphasound-build-time 2>/dev/null)
    [ -n "$fmt" ] || return 0

    ebegin "seeding clock from build time ($fmt)"
    /usr/local/bin/alphasound-checkpoint clock-set
    date -s "$fmt" >/dev/null 2>&1
    eend $?
}
