# overlay/etc/init.d/alphasound-install
#!/sbin/openrc-run
# Alpine diskless mode does not auto-install packages from /etc/apk/world
# at boot — the world file is just a list of desired packages, nothing
# actually calls `apk add`. Without this service, every package we add to
# the build never lands on the running system. We run in the boot runlevel
# so binaries exist before the default runlevel tries to start sshd,
# bluetooth, dnsmasq, etc.

description="Install Alphasound packages from offline cache"

depend() {
    after modules
    keyword -shutdown
}

start() {
    if command -v hostapd >/dev/null 2>&1; then
        einfo "Packages already installed, skipping"
        return 0
    fi
    ebegin "Installing Alphasound packages from offline cache"
    # Install per-package so a single failing dep doesn't block everything
    # — at minimum we want hostapd + dnsmasq + openssh so the AP and ssh
    # come up and the rest can be debugged remotely.
    failed=""
    for pkg in $(cat /etc/apk/world); do
        if ! apk add --no-network --no-progress "$pkg" >/dev/console 2>&1; then
            echo "FAILED to install $pkg" >/dev/console
            failed="$failed $pkg"
        fi
    done
    [ -n "$failed" ] && echo "Some packages failed:$failed" >/dev/console
    eend 0
}
